# services/com.pactoria.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pactoria.frontend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rezazeraat/Desktop/Pactoria-MVP/services/frontend-start.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/rezazeraat/Desktop/Pactoria-MVP/frontend</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>NODE_ENV</key>
        <string>development</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/rezazeraat/Desktop/Pactoria-MVP/services/logs/frontend-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/rezazeraat/Desktop/Pactoria-MVP/services/logs/frontend-stderr.log</string>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>Nice</key>
    <integer>0</integer>
</dict>
</plist>